# crontab
15 22 * * * cd /app && node booking.js >> /var/log/booking.log 2>&1
